# package/rooter/ext-rtorrent/files/etc/init.d/rtorrent
#!/bin/sh /etc/rc.common

START=99

start() {
  sleep 3
  screen -dm -t rtorrent nice -19 rtorrent
}

stop() {
  killall rtorrent
}